# scripts/gspro.service
[Unit]
Description=GSPro Tournament Scoring Appliance
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/gspro
ExecStart=/opt/gspro/start.sh
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
